# conf/bot_sv.conf
[program:GeoBot]
user=root
command=python3 /data/app/GeoBot/start.py
autorstart=true
directory=/data/app/GeoBot
autorestart=true
startsecs=5
startretries=100
redirect_stderr=true